# supervisord.conf
[program:node]
command = node server.js cloud9
directory = /home/dotcloud/current